# systemd/layout_switcher.service
# Service файл для запуска только при активной сессии
[Unit]
Description=Layout Switcher Service 
After=graphical-session.target
Requires=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/projects/my/layout_switcher/start.sh
WorkingDirectory=%h/projects/my/layout_switcher
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus

# Важно: НЕ использовать RemainAfterExit
# Сервис будет остановлен при завершении сессии

[Install]
WantedBy=graphical-session.target